# external/rootfs-overlay/etc/systemd/system/disk-expand-root.service
[Unit]
Description=Extend root partition and resize ext4 file system
After=local-fs.target
Wants=local-fs.target

[Service]
Environment=ROOT_DISK=/dev/vda
Environment=ROOT_PARTITION=1
ExecStart=/bin/sh -c "/usr/bin/partextend \${ROOT_DISK} \${ROOT_PARTITION} || exit 0"
ExecStop=/bin/sh -c "/sbin/resize2fs \${ROOT_DISK}\${ROOT_PARTITION} || exit 0"
Type=oneshot

[Install]
WantedBy=multi-user.target
