# external/rootfs-overlay/etc/systemd/system/cri-docker.service.d/cni.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cri-dockerd --container-runtime-endpoint fd:// --pod-infra-container-image=k8s.gcr.io/pause:3.7 --network-plugin=cni
